# infra/bounce-sweep.timer
[Unit]
Description=Kör bounce-sweep dagligen (150 kommunpolitiker/dag, täcker alla på 90 dagar)

[Timer]
OnCalendar=*-*-* 10:00:00
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
