# infra/letter-generator.service
[Unit]
Description=Politiker-webapp brevgenerator (Claude AI)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=berduf
ExecStart=/usr/bin/python3 /home/berduf/GitHub/politiker-webapp/infra/letter-generator.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=letter-generator
