# infra/monitor.service
[Unit]
Description=Politiker-webapp nyhetsbevakare
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=berduf
ExecStart=/usr/bin/python3 /home/berduf/GitHub/politiker-webapp/infra/monitor.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=monitor
